# supervisord.conf
[supervisord]
nodaemon=true

[program:bbr]
command=/apps/rinetd -f -c /apps/rinetd.conf raw eth0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ssr]
directory=/apps/shadowsocksr-manyuser/shadowsocks
command=python server.py -p 443 -k %(ENV_SSRPASS)s -m %(ENV_SSRENC)s -O %(ENV_SSRPRO)s -o %(ENV_SSROBFS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=%(ENV_PYTHONPATH)s,PATH=%(ENV_PATH)s
